# infrastructure/services/archive/band.service
[Unit]
Description=Band Full Node

[Service]
Environment=BANDD_P2P_LADDR=tcp://0.0.0.0:2210
Environment=BANDD_RPC_LADDR=tcp://0.0.0.0:2211
Environment=BANDD_NODE=tcp://127.0.0.0:2211
Environment=BANDD_GRPC_ADDRESS=0.0.0.0:2212
Environment=BANDD_API_ADDRESS=tcp://127.0.0.0:2213
Environment=BANDD_GRPC_WEB_ADDRESS=127.0.0.0:2214
Environment=BANDD_P2P_SEEDS="8d42bdcb6cced03e0b67fa3957e4e9c8fd89015a@34.87.86.195:26656,543e0cab9c3016a0e99775443a17bcf163038912@34.150.156.78:26656"
Environment=BANDD_P2P_PERSISTENT_PEERS="98823087b61d442a4ab86998709c77b2e517ee78@35.240.152.216:26656,3ea84babead3d6bc488810a0f2cf0744cf5c68fe@34.86.22.251:26656,48f26a2c8d238b01124d49818c3bddba826c4356@136.243.5.247:58856,d047cfabdfd5e244af530d6d2101d07c45ff7424@165.22.167.234:41656,1a4af7cbd3db94a3881dc35cfa261ec2ac788f8f@91.246.64.247:26656,c1b93580023f12891683431e25fb05d6e4372bef@188.34.207.243:26656,2f62588be3cf41a019289097c2068c4e65aa5e37@95.217.213.135:26656,58ce8a2bd441057c0b18d0812f692b2d8af55af7@168.119.61.205:26656,570787c6484fb5aef9182f032dbd54042d93b93c@35.82.85.220:26656,cb777d12ec515884f9f8153e2bbe1844079323a4@44.234.200.116:26656,cacb630ef586ffefccf103d442c02bdbc187f996@3.37.82.201:26656,39d45dae55f36db42ef3997376efbbf725666f75@51.38.53.4:30656,70c542ec2e1aebb385f7e93d1c55f1a4073cdcb7@54.38.46.124:30656,d7e04f00e2701cb9b2c62a78b2e7bb3cdb165dfb@51.91.67.48:30656,96b1e7c754ef5a8f80be569a30de21f9675589c9@13.232.162.252:26656"
Environment=BANDD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=BANDD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=BANDD_PRUNING=nothing
Environment=BANDD_HOME=/osmo2/band
User=root
Restart=always
RestartSec=30
ExecStart=/root/go/bin/bandd start

[Install]
WantedBy=multi-user.target
