# infrastructure/services/archive/ixo.service
[Unit]
Description=IXO Full Node

[Service]
Environment=IXOD_P2P_LADDR=tcp://0.0.0.0:2130
Environment=IXOD_RPC_LADDR=tcp://0.0.0.0:2131
Environment=IXOD_NODE=tcp://127.0.0.1:2131
Environment=IXOD_GRPC_ADDRESS=0.0.0.0:2132
Environment=IXOD_API_ADDRESS=tcp://127.0.0.1:2133
Environment=IXOD_GRPC_WEB_ADDRESS=127.0.0.1:2134
Environment=IXOD_P2P_SEEDS="cbe8c6a5a77f861db8edb1426b734f2cf1fa4020@18.166.133.210:26656,36e4738c7efcf353d3048e5e6073406d045bae9d@80.64.208.42:26656,f0d4546fa5e0c2d84a4244def186b9da3c12ba1a@46.166.138.214:26656,c95af93f0386f8e19e65997262c9f874d1901dc5@18.163.242.188:26656"
Environment=IXOD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=IXOD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=IXOD_PRUNING=nothing
LimitNOFILE=500000
User=root
Restart=always
RestartSec=30
ExecStart=/root/go/bin/ixod start

[Install]
WantedBy=multi-user.target
