# infrastructure/services/archive/kichain.service
[Unit]
Description=Kichain Full Node

[Service]
Environment=KID_P2P_LADDR=tcp://0.0.0.0:2100
Environment=KID_RPC_LADDR=tcp://0.0.0.0:2101
Environment=KID_NODE=tcp://127.0.0.0:2101
Environment=KID_GRPC_ADDRESS=0.0.0.0:2102
Environment=KID_API_ADDRESS=tcp://127.0.0.0:2103
Environment=KID_GRPC_WEB_ADDRESS=127.0.0.0:2104
Environmant=KID_P2P_SEEDS=24cbccfa8813accd0ebdb09e7cdb54cff2e8fcd9@51.89.166.197:26656
Environment=KID_P2P_PERSISTENT_PEERS=81396d4703a2e3cbd136c7324e4df5686fd48218@35.180.8.214:26656,c597db55d9a609b8b77c3d37ecf1fa9a67117cc0@144.217.82.4:26656,50af457247b59aa558a26a14bd7ac4bf86eeae70@195.201.164.223:26656,81eef39d2ca9a07490857d197423da4ba5e01879@15.188.134.35:26656,5adb5ad6a6fcef624866cefdb551dafdc07f7e78@15.188.198.188:26656,41b321292cbe50c5c30017cc71c404481be0e20b@3.38.12.5:26656,7835832e5ded6b0b3db4b1eb1d202adfe1c67629@161.97.132.109:26656
Environment=KID_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=KID_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=KID_PRUNING=nothing
LimitNOFILE=500000
User=root
Restart=always
RestartSec=30
ExecStart=/root/go/bin/kid start

[Install]
WantedBy=multi-user.target
