# infrastructure/services/archive/osmosis.service
[Unit]
Description=Osmosis Full Node

[Service]
Environment=OSMOSISD_P2P_LADDR=tcp://0.0.0.0:2000
Environment=OSMOSISD_RPC_LADDR=tcp://0.0.0.0:2001
Environment=OSMOSISD_NODE=tcp://127.0.0.1:2001
Environment=OSMOSISD_GRPC_ADDRESS=0.0.0.0:2002
Environment=OSMOSISD_API_ADDRESS=tcp://127.0.0.1:2003
Environment=OSMOSISD_GRPC_WEB_ADDRESS=127.0.0.1:2004
Environment=OSMOSISD_P2P_UNCONDITIONAL_PEERS=d22a0878f4694d7d952e3e7053fc56a8840a520b,817d66de08ebd3f81ef434e9d7188092f63768b4,28f61c154c82f0122a841a12f8aa87703bd6ae1e,d49f8df41c728709740d7012b671346287d82cf7
Environment=OSMOSISD_P2P_SEED_NODES=1b077d96ceeba7ef503fb048f343a538b2dcdf1b@136.243.218.244,2308bed9e096a8b96d2aa343acc1147813c59ed2@3.225.38.25:26656,902bdfe51b6a97cc9369664a21c87ed61d471d2a@136.243.218.243,085f62d67bbf9c501e8ac84d4533440a1eef6c45@95.217.196.54:26656,f515a8599b40f0e84dfad935ba414674ab11a668@osmosis.blockpane.com:26656
Environment=OSMOSISD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=OSMOSISD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=OSMOSISD_PRUNING=nothing
Environment=OSMOSISD_HOME=/osmo1/
LimitNOFILE=500000
User=root
Restart=always
RestartSec=30
ExecStart=/root/go/bin/osmosisd start

[Install]
WantedBy=multi-user.target
