# infrastructure/services/archive/sentinel.service
 
[Unit]
Description=Sentinel Hub

[Service]
Environment=SENTINELHUB_P2P_LADDR=tcp://0.0.0.0:2030
Environment=SENTINELHUB_RPC_LADDR=tcp://0.0.0.0:2031
Environment=SENTINELHUB_GRPC_ADDRESS=127.0.0.1:2032
Environment=SENTINELHUB_API_ADDRESS=tcp://127.0.0.1:2033
Environment=SENTINELHUB_NODE=tcp://127.0.0.1:2031
Environment=SENTINELHUB_P2P_PERSISTENT_PEERS="152bad169bfde238d24b8e2403d72c092adc107b@54.169.185.248:26656,57dbef4a8637f4c14dc5f8a4070ff07cc3de8380@13.212.44.171:26656,cb9a308bd21b745ecd58236a08bc60046b06fec3@13.233.250.145:26656,a77f6a094578dad899e2f40e0626b4c6d4705311@3.36.165.232:26656,5f7164bbf0e74bc0b7bc624c631b559bde6d52b4@165.22.65.48:26656,647e54e22cfcdf97b42ada11fcef2f0b5bd81230@35.163.249.152:31358,387027e3b1180d3a619cbbf3462704a490785963@54.176.90.228:26656"
Environment=SENTINELHUB_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=SENTINELHUB_STATE_SYNC_SNAPSHOT_INTERVAL=1000
Environment=HOME=/root
LimitNOFILE=500000
ExecStart=/usr/bin/sentinelhub start --pruning nothing
LogLevelMax=3

[Install]
WantedBy=multi-user.target
