# infrastructure/services/gaia.service
[Unit]
Description=Gaia

[Service]
Environment=GAIAD_P2P_LADDR=tcp://0.0.0.0:2010
Environment=GAIAD_RPC_LADDR=tcp://0.0.0.0:2011
Environment=GAIAD_GRPC_ADDRESS=127.0.0.1:2012
Environment=GAIAD_API_ADDRESS=tcp://127.0.0.1:2013
Environment=GAIAD_NODE=tcp://127.0.0.1:2011
Environment=GAIAD_P2P_PERSISTENT_PEERS="bf8328b66dceb4987e5cd94430af66045e59899f@public-seed.cosmos.vitwit.com:26656,cfd785a4224c7940e9a10f6c1ab24c343e923bec@164.68.107.188:26656,d72b3011ed46d783e369fdf8ae2055b99a1e5074@173.249.50.25:26656,ba3bacc714817218562f743178228f23678b2873@public-seed-node.cosmoshub.certus.one:26656,3c7cad4154967a294b3ba1cc752e40e8779640ad@84.201.128.115:26656,366ac852255c3ac8de17e11ae9ec814b8c68bddb@51.15.94.196:26656"
Environment=GAIAD_P2P_UNCONDITIONAL_PEERS="11bfced94d58357917fe9189fb34d3def790b454"
Environment=GAIAD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=GAIAD_SNAPSHOT_INTERVAL=1000
Environment=GAIAD_P2P_MAX_NUM_OUTBOUND_PEERS=60
Environment=HOME=/root/
LimitNOFILE=500000
ExecStart=/usr/bin/gaiad start --pruning custom --pruning-keep-recent 362880 --pruning-keep-every 0 --pruning-interval 100
LogLevelMax=3


[Install]
WantedBy=multi-user.target
