# infrastructure/services/validator/emoney.service
[Unit]
Description=E-Money Validator

[Service]
Environment=EMD_P2P_LADDR=tcp://0.0.0.0:2140
Environment=EMD_RPC_LADDR=tcp://0.0.0.0:2141
Environment=EMD_GRPC_ADDRESS=127.0.0.1:2142
Environment=EMD_API_ADDRESS=tcp://127.0.0.1:2143
Environment=EMD_GRPC_WEB_ADDRESS=127.0.0.1:2144
Environment=EMD_P2P_SEEDS="6420ef5087accdff4a87df5331d07da5de568743@18.194.208.47:28656,f49bf0e3d6d6057499ceb6613854af37a3da532a@3.121.126.177:28656,ecec8933d80da5fccda6bdd72befe7e064279fc1@207.180.213.123:26676,0ad7bc7687112e212bac404670aa24cd6116d097@50.18.83.75:26656,1723e34f45f54584f44d193ce9fd9c65271ca0b3@13.124.62.83:26656,34eca4a9142bf9c087a987b572c114dad67a8cc5@172.105.148.191:26656,0b186517e4d82eb4c000a567e486b7b96bf19752@44.195.95.22:26656,eed66085c975189e3d498fe61af2fcfb3da34924@217.79.184.40:26656,226d713acb10a1cd09b3b890406b886d7b071cf2@10.105.2.101:26656"
Environment=EMD_STATE_SYNC_SNAPSHOT_INTERVAL=1000
Environment=EMD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=EMD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=EMD_PRUNING=custom
Environment=EMD_PRUNING_KEEP_RECENT=100
Environment=EMD_PRUNING_KEEP_EVERY=0
Environment=EMD_PRUNING_INTERVAL=100
Environment=EMD_CONSENSUS_DOUBLE_SIGN_CHECK_HEIGHT=10
Environment=EMD_TX_INDEX_INDEXER=null
Environment=EMD_HOME=/chains/em
User=cyber
Restart=always
RestartSec=30
LimitNOFILE=500000
ExecStart=/home/cyber/go/bin/emd start

[Install]
WantedBy=multi-user.target
