# infrastructure/services/validator/juno.service
[Unit]
Description=Juno

[Service]
Environment=JUNOD_P2P_LADDR=tcp://0.0.0.0:2070
Environment=JUNOD_RPC_LADDR=tcp://0.0.0.0:2071
Environment=JUNOD_GRPC_ADDRESS=127.0.0.1:2072
Environment=JUNOD_API_ADDRESS=tcp://127.0.0.1:2073
Environment=JUNOD_NODE=tcp://127.0.0.1:2071
Environment=JUNOD_P2P_SEEDS="b1f46f1a1955fc773d3b73180179b0e0a07adce1@162.55.244.250:39656,7f593757c0cde8972ce929381d8ac8e446837811@178.18.255.244:26656,7b22dfc605989d66b89d2dfe118d799ea5abc2f0@167.99.210.65:26656,4bd9cac019775047d27f9b9cea66b25270ab497d@137.184.7.164:26656,bd822a8057902fbc80fd9135e335f0dfefa32342@65.21.202.159:38656,15827c6c13f919e4d9c11bcca23dff4e3e79b1b8@51.38.52.210:38656,e665df28999b2b7b40cff2fe4030682c380bf294@188.40.106.109:38656,92804ce50c85ff4c7cf149d347dd880fc3735bf4@34.94.231.154:26656,795ed214b8354e8468f46d1bbbf6e128a88fe3bd@34.127.19.222:26656,dc2540dabadb8302da988c95a3c872191061aed2@p2p.7d1b53c0-b86b-44c8-8c02-e3b0e88a4bf7.cryptodotorg.herd.run:26656,33b15c14f54f71a4a923ac264761eb3209784cf2@p2p.0d20d4b3-6890-4f00-b9f3-596ad3df6533.cryptodotorg.herd.run:26656,d2862ef8f86f9976daa0c6f59455b2b1452dc53b@p2p.a088961f-5dfd-4007-a15c-3a706d4be2c0.cryptodotorg.herd.run:26656,87c3adb7d8f649c51eebe0d3335d8f9e28c362f2@seed-0.crypto.org:26656,e1d7ff02b78044795371beb1cd5fb803f9389256@seed-1.crypto.org:26656,2c55809558a4e491e9995962e10c026eb9014655@seed-2.crypto.org:26656"
Environment=JUNOD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=JUNOD_STATE_SYNC_SNAPSHOT_INTERVAL=1000
Environment=JUNOD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=JUNOD_PRUNING=custom
Environment=JUNOD_PRUNING_KEEP_RECENT=100
Environment=JUNOD_PRUNING_KEEP_EVERY=0
Environment=JUNOD_PRUNING_INTERVAL=100
Environment=JUNOD_CONSENSUS_DOUBLE_SIGN_CHECK_HEIGHT=10
Environment=JUNOD_TX_INDEX_INDEXER=null
Environment=JUNOD_HOME=/chains/juno-mainnet
User=cyber
Restart=always
RestartSec=30
LimitNOFILE=500000
ExecStart=/home/cyber/go/bin/junod start

[Install]
WantedBy=multi-user.target
