# infrastructure/services/validator/stargaze.service
[Unit]
Description=Stargaze Validator

[Service]
Environment=STARSD_P2P_LADDR=tcp://0.0.0.0:4000
Environment=STARSD_RPC_LADDR=tcp://0.0.0.0:4001
Environment=STARSD_GRPC_ADDRESS=127.0.0.1:4002
Environment=STARSD_API_ADDRESS=tcp://127.0.0.1:4003
Environment=STARSD_GRPC_WEB_ADDRESS=127.0.0.1:4004
Environment=STARSD_P2P_SEEDS="a54a9d04a8b42a0bd495a93f18efcf2a272277e0@162.55.132.230:4000,1d73521c565b37a53038fc730bcd207a3db361b6@144.91.91.30:26656,320e4b81ab327dd2593a39de0d3ae718fdb9347c@176.9.168.220:26656,413562a502b4a47ec3375bfae6aa01bc47bc6979@204.236.141.188:26656,e28635dc537dce1d6f1cf44b7b335a530d923fa1@15.165.94.246:26656,387c32677c54dd4627366ffdf8e2dda68e71af69@54.193.168.147:26656,1de5e685e687d0df35e6f3734ee04aff956d5c0f@75.119.154.98:26656"
Environment=STARSD_STATE_SYNC_SNAPSHOT_INTERVAL=1000
Environment=STARSD_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=STARSD_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=STARSD_PRUNING=custom
Environment=STARSD_PRUNING_KEEP_RECENT=100
Environment=STARSD_PRUNING_KEEP_EVERY=0
Environment=STARSD_PRUNING_INTERVAL=100
Environment=STARSD_CONSENSUS_DOUBLE_SIGN_CHECK_HEIGHT=10
Environment=STARSD_TX_INDEX_INDEXER=null
Environment=STARSD_HOME=/chains/.starsd
User=cyber
Restart=always
RestartSec=30
LimitNOFILE=500000
ExecStart=/home/cyber/go/bin/starsd start

[Install]
WantedBy=multi-user.target
