# infrastructure/services/validator/starname.service
[Unit]
Description=Starname Validator

[Service]
Environment=STARNAMED_P2P_LADDR=tcp://0.0.0.0:2170
Environment=STARNAMED_RPC_LADDR=tcp://0.0.0.0:2171
Environment=STARNAMED_GRPC_ADDRESS=127.0.0.1:2172
Environment=STARNAMED_API_ADDRESS=tcp://127.0.0.1:2173
Environment=STARNAMED_GRPC_WEB_ADDRESS=127.0.0.1:2174
Environment=STARNAMED_P2P_SEEDS="ca133187b37b59d2454812cfcf31b6211395adec@167.99.194.126:16656,1c7e014b65f7a3ea2cf48bffce78f5cbcad2a0b7@13.37.85.253:26656,8c64a2127cc07d4570756b61f83af60d34258398@13.37.61.32:26656,9aabe0ac122f3104d8fc098e19c66714c6f1ace9@3.37.140.5:26656,faedef1969911d24bf72c56fc01326eb891fa3b7@63.250.53.45:16656,94ac1c02b4e2ca3fb2706c91a68b8030ed3615a1@35.247.175.128:16656,be2235996b1c785a9f57eed25fd673ca111f0bae@52.52.89.64:26656,f63d15ab7ed55dc75f332d0b0d2b01d529d5cbcd@212.71.247.11:26656,f5597a7ed33bc99eb6ba7253eb8ac76af27b4c6d@138.201.20.147:26656"
Environment=STARNAMED_STATE_SYNC_SNAPSHOT_INTERVAL=1000
Environment=STARNAMED_P2P_MAX_NUM_INBOUND_PEERS=500
Environment=STARNAMED_P2P_MAX_NUM_OUTBOUND_PEERS=100
Environment=STARNAMED_PRUNING=custom
Environment=STARNAMED_PRUNING_KEEP_RECENT=100
Environment=STARNAMED_PRUNING_KEEP_EVERY=0
Environment=STARNAMED_PRUNING_INTERVAL=100
Environment=STARNAMED_CONSENSUS_DOUBLE_SIGN_CHECK_HEIGHT=10
Environment=STARNAMED_TX_INDEX_INDEXER=null
Environment=STARNAMED_HOME=/spin/starnamedata/starnamedata
User=cyber
Restart=always
RestartSec=30
LimitNOFILE=500000
ExecStart=/home/cyber/go/bin/starnamed start

[Install]
WantedBy=multi-user.target
